# supervisord.conf
[supervisord]
nodaemon=true

[program:graphql-frontend]
command=/app/start-service.sh